# Translation_System/configs/etc/udev/rules.d/00-audio_off_poweroff.rules
ACTION=="remove", KERNEL=="pcmC1D0p", RUN+="/bin/systemctl poweroff"
#ACTION=="remove", ATTRS{idVendor}=="047f", RUN+="/bin/systemctl poweroff"
