# Translation_System/configs/etc/udev/rules.d/01-USBKEY_copy.rules
ACTION=="add", DRIVERS=="usb-storage", KERNEL=="sda1", RUN+="/usr/local/sbin/copy_recording.sh"
